# miAplicacion.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid ramon
chdir /var/www/public_html/

env PORT=8080

exec npm start